# auto-big-picture.service.template
[Unit]
Description=Handles steam big picture for Bluetooth/USB controllers

[Service]
ExecStart=__HOME__/scripts/auto_big_picture.py
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
