# auto-big-picture.service.template
[Unit]
Description=Handles steam big picture for Bluetooth/USB controllers

[Service]
ExecStart=__SCRIPT_PATH__
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
